# _infra/dev/supervisord.conf
[supervisord]
nodaemon=true

[program:ng_build_watch]
command=npm run build -- --watch --configuration=development
directory=/srv/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ng_build_watch.log
stderr_logfile=/var/log/supervisor/ng_build_watch_err.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
